# systemd/conclave-audit.service
[Unit]
Description=Conclave RAG/MCP audit
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/conclave
ExecStart=/usr/bin/python3 -m conclave.cli audit --mode all --output-dir %h/.conclave/audits
