# systemd/conclave-validate.timer
[Unit]
Description=Conclave nightly validation timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
