# systemd/conclave-audit.timer
[Unit]
Description=Conclave audit weekly timer

[Timer]
OnCalendar=Sun *-*-* 05:15:00
Persistent=true

[Install]
WantedBy=timers.target
