# systemd/conclave-reconcile.service
[Unit]
Description=Conclave weekly reconcile
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/conclave
ExecStart=/usr/bin/python3 -m conclave.cli reconcile --topic all
